# debian/prerm
#!/bin/sh
set -e
systemctl stop    picam-orchestrator 2>/dev/null || true
systemctl disable picam-orchestrator 2>/dev/null || true
